# configs/etc/service/exp90/run
#!/bin/bash
HOME="/home/exp90"
USER="exp90"
cd $HOME
exec 2>&1
exec setuidgid exp90 /usr/local/bin/tcpserver -H -R -c 500 -u exp90 0 12589 /opt/nodejs/bin/node task.js 
