# configs/etc/service/code70/run
#!/bin/bash
HOME="/home/code70"
USER="code70"
cd $HOME
exec 2>&1
exec setuidgid code70 python2 code70.py 
