# configs/etc/service/code80/run
#!/bin/bash
HOME="/home/code80"
USER="code80"
cd $HOME
exec 2>&1
exec setuidgid code80 python2 code80.py 
